# backend.service
[Unit]
Description=Mongodb

[Service]
User=mongodb

WorkingDirectory=/home/mongodb
ExecStart=node ./Index.js

[Install]
WantedBy=multi-user.target
