# mongodbadm.service
[Unit]
Description=Mongodb

[Service]
User=mongodbadm

WorkingDirectory=/home/mongodbadm


[Install]
WantedBy=multi-user.target
